# docker/s6-rc.d/api/run
#!/bin/sh
exec node /opt/meta-server.js
