# docker/s6-rc.d/api/run
#!/bin/sh
exec /command/with-contenv node /opt/meta-server.js
