# docker/s6-rc.d/init-config/up
#!/bin/sh
exec /usr/bin/with-contenv /opt/generate-nginx-conf.sh
